# netfb.socket
[Unit]
Description=netfb socket

[Socket]
ListenStream=11235
Accept=yes

[Install]
WantedBy=sockets.target
